# install/home/pi/crontab
@reboot cd /home/pi/netmonitor && launch.sh

